# docker/production/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:scheduler]
directory=/var/www/html
command=php artisan schedule:work
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:queue-calls]
directory=/var/www/html
command=php artisan queue:work redis --queue=calls --tries=1 --timeout=0 --sleep=1
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:queue-retry]
directory=/var/www/html
command=php artisan queue:work redis --queue=calls-retry --tries=1 --timeout=0 --sleep=1
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:kafka-incoming]
directory=/var/www/html
command=/bin/sh -lc 'while true; do php artisan calls:kafka:consume "${KAFKA_INCOMING_CALLS_TOPIC:-incoming-calls}" --group="${KAFKA_INCOMING_CALLS_GROUP:-calls-incoming}" --source=incoming-calls-consumer --limit="${KAFKA_CONSUMER_BATCH_LIMIT:-1000}" --timeout-ms="${KAFKA_CONSUMER_TIMEOUT_MS:-1000}"; done'
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:kafka-telephony-facts]
directory=/var/www/html
command=/bin/sh -lc 'while true; do php artisan calls:kafka:consume "${KAFKA_TELEPHONY_FACTS_TOPIC:-telephony.facts}" --group="${KAFKA_TELEPHONY_FACTS_GROUP:-calls-telephony-facts}" --source=telephony-facts-consumer --limit="${KAFKA_CONSUMER_BATCH_LIMIT:-1000}" --timeout-ms="${KAFKA_CONSUMER_TIMEOUT_MS:-1000}"; done'
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:outbox-publisher]
directory=/var/www/html
command=/bin/sh -lc 'while true; do php artisan calls:telephony-outbox:publish --limit="${TELEPHONY_OUTBOX_PUBLISH_LIMIT:-100}"; sleep "${OUTBOX_PUBLISH_INTERVAL_SECONDS:-1}"; done'
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
